# pi-backup/com.vybhavab.photos-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vybhavab.photos-backup</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/vybhavab/.local/bin/backup-photos</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/vybhavab/.local/share/pi-backup/launchd-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/vybhavab/.local/share/pi-backup/launchd-stderr.log</string>
</dict>
</plist>
