# xppen-ack05/com.vybhavab.ack05.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vybhavab.ack05</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/vybhavab/.bun/bin/bun</string>
        <string>run</string>
        <string>/Users/vybhavab/projects/xptool/index.ts</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/ack05.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/ack05.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
